# apache2.conf
[program:apache2]
directory=/usr/sbin/
command=/usr/sbin/apache2 -D FOREGROUND
environment=APACHE_RUN_USER='www-data',APACHE_RUN_GROUP='www-data',APACHE_LOG_DIR='/var/log/apache2'
autostart=true
autorestart=true
redirect_stderr=true
