# dataloop.conf
[program:dataloop]
directory=/usr/local/bin/
command=/usr/local/bin/dataloop-lin-agent --api-key a59cd326-4ff6-4dab-a93c-d182c8456a58 --tags=web
autostart=true
autorestart=true
redirect_stderr=true
